# shadowsocks/files/rc.local
#!/bin/bash
touch /var/lock/subsys/local
ssserver  -c  /etc/shadowsocks.json -d  start --user  ssuser
